# ops/launchd/rocks.shine.dfarm-infra.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>rocks.shine.dfarm-infra</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd ~/projects/sorrentolab/device-farm &amp;&amp; mise run boot:infra</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/tmp/dfarm-infra.log</string>
  <key>StandardErrorPath</key><string>/tmp/dfarm-infra.log</string>
</dict>
</plist>
